# stage2/03-docker/files/random-reboot.service
[Unit]
Description=Wait a random period then restart

[Service]
ExecStart=python /usr/sbin/random-reboot.py

